# lxlauncher-vala/etc/xdg/autostart/lxlauncher.desktop
[Desktop Entry]
Type=Application
Name=LxLauncher
Exec=lxlauncher
NoDisplay=true

